# roles/niflheim_profile/files/profile.d/omp_num_threads.sh
# OpenMP default number of threads to use in parallel regions
export OMP_NUM_THREADS=1
